# gomycode/github/GENREAL/systemd/backup-demo.service
[Unit]
Description=Run demo backup script

[Service]
Type=oneshot
ExecStart=/opt/devops-starter/scripts/backup-demo.sh

